# okawhisp.service
[Unit]
Description=OkaWhisp - System-Level Voice Input (F9 Hotkey)
Documentation=https://github.com/YOUR_USER/okawhisp
After=graphical-session.target pipewire.service
StartLimitBurst=3
StartLimitIntervalSec=60s

[Service]
Type=simple
# NOTE: This is a template for local dev.
# The installer (install.sh) generates the actual service file with correct paths.
ExecStart=/usr/bin/env bash -c 'uv run ~/.local/share/okawhisp/okawhisp.py --key F9 --model large-v3 --engine faster --language de'

# X11 + Audio Zugriff
Environment="DISPLAY=:1"
Environment="XAUTHORITY=/home/okahari/.Xauthority"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="PYTHONUNBUFFERED=1"
Environment="CUDA_VISIBLE_DEVICES=0"

Restart=always
RestartSec=5s

# Ressourcen-Limits
MemoryMax=4G
MemoryHigh=3G

StandardOutput=journal
StandardError=journal
SyslogIdentifier=okawhisp

[Install]
WantedBy=graphical-session.target
