# voice-type.service.example
[Unit]
Description=Voice Type - System-Level Voice Input (F9 Hotkey)
After=graphical-session.target pipewire.service
StartLimitBurst=3
StartLimitIntervalSec=60s

[Service]
Type=simple
# uv manages dependencies automatically — no venv setup needed
ExecStart=%h/.local/bin/uv run %h/.local/share/voice-type/voice-type.py \
    --key F9 \
    --model large-v3 \
    --engine faster \
    --language de \
    --beam-size 5 \
    --silence 2.0

# Adjust these to match your display setup (run: echo $DISPLAY $XAUTHORITY)
Environment="DISPLAY=:0"
Environment="XAUTHORITY=%h/.Xauthority"
Environment="XDG_RUNTIME_DIR=/run/user/%i"
Environment="PYTHONUNBUFFERED=1"
Environment="CUDA_VISIBLE_DEVICES=0"

Restart=always
RestartSec=5s
MemoryMax=4G
MemoryHigh=3G

StandardOutput=journal
StandardError=journal
SyslogIdentifier=voice-type

[Install]
WantedBy=graphical-session.target
